# helmet/deploy/systemd/helmet-voice.service
[Unit]
Description=Helmet Voice Service
Documentation=https://github.com/yourusername/helmet
After=network.target sound.target
Requires=network.target sound.target helmet-orchestrator.service

[Service]
Type=simple
User=helmet
Group=helmet
WorkingDirectory=/opt/helmet
Environment=HELMET_PROFILE=field
Environment=PYTHONPATH=/opt/helmet:/opt/helmet/libs
Environment=PULSE_RUNTIME_PATH=/run/user/1000/pulse
ExecStart=/opt/helmet/venv/bin/python /opt/helmet/services/voice/voice_service.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helmet-voice

# Resource limits
MemoryMax=1G
CPUQuota=60%

# Audio access
SupplementaryGroups=audio pulse-access

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/helmet/logs /opt/helmet/models
PrivateTmp=true

[Install]
WantedBy=helmet.target